# setup/supervisor-wayloshare-worker.conf
# Supervisor Configuration for WayloShare Queue Workers
# Location: /etc/supervisor/conf.d/wayloshare-worker.conf
# After creating this file, run:
#   sudo supervisorctl reread
#   sudo supervisorctl update
#   sudo supervisorctl start wayloshare-worker:*

[program:wayloshare-worker]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/wayloshare/artisan queue:work redis --sleep=3 --tries=3 --max-time=3600
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=4
redirect_stderr=true
stdout_logfile=/var/www/wayloshare/storage/logs/worker.log
stopwaitsecs=3600

# Logging
stderr_logfile=/var/www/wayloshare/storage/logs/worker_error.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
